# config/production/upstart/instance.conf
description     "One Block Together service, started by blocktogether.conf"

start on (local-filesystems and net-device-up)
stop on shutdown

setuid ubuntu
setgid ubuntu
chdir /usr/local/blocktogether/current

instance "Block Together - ${NAME}"

# Give up if respawn occurs 10 times in 10 seconds.
respawn
respawn limit 10 10

pre-start script
    ulimit -n 64000
end script

exec /usr/bin/js ${NAME}.js >> /data/logs/${NAME}.log
